# _misc/arch_package/hadoop/hadoop-secondarynamenode.service
[Unit]
Description=Haoop secondary name node server

[Service]
Type=forking
User=hadoop
Group=hadoop
EnvironmentFile=/etc/conf.d/hadoop
ExecStart=/usr/bin/hdfs --daemon start secondarynamenode
ExecStop=/usr/bin/hdfs --daemon stop secondarynamenode
Restart=on-abnormal
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
